# infra/systemd/ha-daily-ac-schedule-midnight-heat.timer
[Unit]
Description=Daily AC Schedule Timer - Midnight Heat

[Timer]
OnCalendar=*-*-* 23:50:00
Persistent=true
Unit=ha-daily-ac-schedule@midnight-heat.service

[Install]
WantedBy=timers.target
