# infra/systemd/ha-daily-ac-schedule@.service
[Unit]
Description=Daily AC Schedule Slot (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/architect/matrix
ExecStart=/usr/bin/env bash /home/architect/matrix/ops/ha/run_daily_ac_schedule.sh --slot %i
