# infra/systemd/ha-daily-ac-schedule-morning-heat.timer
[Unit]
Description=Daily AC Schedule Timer - Morning Heat

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=true
Unit=ha-daily-ac-schedule@morning-heat.service

[Install]
WantedBy=timers.target
