# infra/systemd/oracle-signal-bridge.service
[Unit]
Description=Oracle Signal Bridge
After=network-online.target signal-oracle-bridge.service
Wants=network-online.target signal-oracle-bridge.service

[Service]
Type=simple
User=oracle
Group=oracle
WorkingDirectory=/home/oracle/oraclebot
Environment=PYTHONUNBUFFERED=1
Environment=TELEGRAM_RUNTIME_ROOT=/home/oracle/oraclebot
Environment=TELEGRAM_SHARED_CORE_ROOT=/home/architect/matrix
Environment=TELEGRAM_EXECUTOR_CMD=/home/oracle/oraclebot/src/telegram_bridge/executor.sh
EnvironmentFile=-/etc/default/oracle-signal-bridge
ExecStartPre=/usr/bin/python3 /home/oracle/oraclebot/src/telegram_bridge/wait_for_signal_transport.py
ExecStart=/usr/bin/python3 /home/oracle/oraclebot/src/telegram_bridge/main.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
