# infra/systemd/server3-runtime-observer.timer
[Unit]
Description=Server3 Runtime KPI Observer Timer

[Timer]
OnCalendar=*-*-* 08:05:00
RandomizedDelaySec=30
Persistent=true
Unit=server3-runtime-observer.service

[Install]
WantedBy=timers.target
