# infra/systemd/server3-ralph-loop.timer
[Unit]
Description=Server3 Ralph Optimization Loop Timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true
Unit=server3-ralph-loop.service

[Install]
WantedBy=timers.target
