# infra/systemd/telegram-architect-memory-restore-drill.service
[Unit]
Description=Telegram Architect Memory Restore Drill
After=network-online.target
Wants=network-online.target
OnFailure=telegram-architect-memory-alert@%n.service

[Service]
Type=oneshot
User=architect
Group=architect
WorkingDirectory=/home/architect/matrix
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/telegram-architect-bridge
ExecStart=/bin/bash /home/architect/matrix/ops/telegram-bridge/memory_restore_drill.sh
