# infra/systemd/telegram-architect-memory-health.service
[Unit]
Description=Telegram Architect Memory Health Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/home/architect/matrix
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/telegram-architect-bridge
ExecStart=/bin/bash /home/architect/matrix/ops/telegram-bridge/memory_health_check.sh
