# infra/systemd/telegram-architect-memory-health.timer
[Unit]
Description=Telegram Architect Memory Health Timer

[Timer]
OnCalendar=*:0/30
RandomizedDelaySec=120
Persistent=true
Unit=telegram-architect-memory-health.service

[Install]
WantedBy=timers.target
