# infra/systemd/telegram-diary-bridge.service
[Unit]
Description=Telegram Diary Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=diary
Group=diary
WorkingDirectory=/home/diary/diarybot
Environment=PYTHONUNBUFFERED=1
Environment=UNIT_NAME=telegram-diary-bridge.service
Environment=TELEGRAM_RUNTIME_ROOT=/home/diary/diarybot
Environment=TELEGRAM_SHARED_CORE_ROOT=/home/architect/matrix
Environment=TELEGRAM_EXECUTOR_CMD=/home/diary/diarybot/src/telegram_bridge/executor.sh
EnvironmentFile=-/etc/default/telegram-diary-bridge
ExecStart=/usr/bin/python3 /home/diary/diarybot/src/telegram_bridge/main.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
