# infra/systemd/telegram-helper-bridge.service
[Unit]
Description=Telegram Helper Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helperbot
Group=helperbot
WorkingDirectory=/home/helperbot/helperbot
Environment=PYTHONUNBUFFERED=1
Environment=UNIT_NAME=telegram-helper-bridge.service
Environment=TELEGRAM_EXECUTOR_CMD=/home/helperbot/helperbot/src/telegram_bridge/executor_helper.sh
EnvironmentFile=-/etc/default/telegram-helper-bridge
ExecStart=/usr/bin/python3 /home/helperbot/helperbot/src/telegram_bridge/main.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/helperbot

[Install]
WantedBy=multi-user.target
