# infra/systemd/telegram-trinity-bridge.service
[Unit]
Description=Telegram Trinity Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trinity
Group=trinity
WorkingDirectory=/home/trinity/trinitybot
Environment=PYTHONUNBUFFERED=1
Environment=UNIT_NAME=telegram-trinity-bridge.service
Environment=TELEGRAM_RUNTIME_ROOT=/home/trinity/trinitybot
Environment=TELEGRAM_EXECUTOR_CMD=/home/trinity/trinitybot/src/telegram_bridge/executor.sh
EnvironmentFile=-/etc/default/telegram-trinity-bridge
ExecStart=/usr/bin/python3 /home/trinity/trinitybot/src/telegram_bridge/main.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
