# infra/system/tv-desktop/home-tv/.config/autostart/server3-tv-brave.desktop
[Desktop Entry]
Type=Application
Name=Server3 TV Session
Comment=Launch Brave fullscreen and prefer HDMI audio
Exec=/home/tv/.local/bin/server3-tv-session-start.sh
Terminal=false
X-GNOME-Autostart-enabled=true
